# supervisor/tasks.mimoun1997.scool.cat.LOCAL.conf
[program:laravel-worker-tasks-mimoun1997-scool-cat]
process_name=%(program_name)s_%(process_num)02d
command=php /home/mimoun/Code/mimoun1997/tasks/artisan queue:work redis --sleep=3 --tries=3
autostart=true
autorestart=true
user=mimoun
numprocs=8
redirect_stderr=true
stdout_logfile=/home/mimoun/Code/mimoun1997/tasks/storage/logs/worker.log
